# deploy/bar-gestionale.service
[Unit]
Description=Bar Gestionale — FastAPI backend
After=network.target

[Service]
Type=simple
User=fra
WorkingDirectory=/opt/bar-gestionale
EnvironmentFile=-/opt/bar-gestionale/.env
ExecStart=/opt/bar-gestionale/.venv/bin/python -m backend.main
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
